# deploy/com.devforge.mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>             <string>com.devforge.mcp</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/YOUR_USERNAME/.local/bin/devforge-mcp</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>DEV_FORGE_CONFIG</key>
    <string>/Users/YOUR_USERNAME/.config/devforge/config.json</string>
  </dict>
  <key>RunAtLoad</key>         <true/>
  <key>KeepAlive</key>         <true/>
  <key>StandardOutPath</key>   <string>/tmp/devforge-mcp.log</string>
  <key>StandardErrorPath</key> <string>/tmp/devforge-mcp.err</string>
</dict>
</plist>
