# deploy/devforge-mcp.service
[Unit]
Description=DevForge MCP Server
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/devforge-mcp
Restart=on-failure
RestartSec=5
Environment=DEV_FORGE_CONFIG=%h/.config/devforge/config.json
WorkingDirectory=%h/.local/share/devforge

[Install]
WantedBy=default.target
